# eclipse/rnis/templates/rnis
#!/sbin/openrc-run
  
name=rins
command="/root/run.sh"
command_user="root"
pidfile="/run/rins/rins.pid"
command_background="yes"

depend() {
        need net
}

start_pre() {
        checkpath --directory --owner $command_user:$command_user --mode 0775 \
                /run/rins /var/log/rins
}