# ict18/view/templates/backend
#!/sbin/openrc-run
  
name=view
command="/usr/bin/node /var/www/localhost/htdocs/backend.js"
command_user="root"
pidfile="/run/view/view.pid"
command_background="yes"

depend() {
        need net
}

start_pre() {
        checkpath --directory --owner $command_user:$command_user --mode 0775 \
                /run/view /var/log/view
}